# systemd/user/rclone-bisync.timer
[Unit]
Description=Timer for Rclone Google Drive Bisync

[Timer]
OnBootSec=1min
OnUnitInactiveSec=10min

[Install]
WantedBy=timers.target
